# home/Library/LaunchAgents/homebrew.simonsimcity.update-notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
      <key>EnableGlobbing</key>
      <false/>
      <key>Label</key>
      <string>homebrew.simonsimcity.update-notifier</string>
      <key>ProgramArguments</key>
      <array>
          <string>/bin/bash</string>
          <string>/Users/jackyalcine/.local/bin/homebrew-update-notifier</string>
      </array>
      <key>RunAtLoad</key>
      <true/>
      <key>StandardErrorPath</key>
      <string>/tmp/homebrew.simonsimcity.update-notifier.err</string>
      <key>StandardOutPath</key>
      <string>/tmp/homebrew.simonsimcity.update-notifier.out</string>
      <key>StartCalendarInterval</key>
      <array>
          <dict>
              <key>Hour</key>
              <integer>10</integer>
              <key>Minute</key>
              <integer>0</integer>
          </dict>
          <dict>
              <key>Hour</key>
              <integer>15</integer>
              <key>Minute</key>
              <integer>0</integer>
          </dict>
      </array>
  </dict>
</plist>
